# home/private_dot_config/systemd/user/xdg-desktop-portal.service.d/override.conf
[Service]
Environment="XDG_CURRENT_DESKTOP=i3"
Environment="DISPLAY=:0"
